# systemd/system/mounter.service
[Service]
ExecStart=/home/sur/src/dotfiles/systemd/mounter.sh

[Install]
WantedBy=multi-user.target
